# s6/init/finish
#!/bin/sh
s6-svscanctl -t "${S6_USER_SCAN}"
exit 125
