# system/s6-rc/services/wlan0.wpa_supplicant/run
#!/bin/sh
. ./data/options.conf
exec "$wpa_supplicant" -i wlan0 -c /etc/wpa_supplicant.conf -D bsd
