# system/s6-rc/services/wlan0.netif/up
define child wlan0
define parent bwn0
/sbin/ifconfig wlan create wlandev $parent name $child
