# system/s6/init/finish
#!/bin/sh
rm -r $(realpath /run/live)
rm /run/live
s6-rc -d change default
