# s6-rc/services/powerdxx/run
#!/bin/sh
. ./data/options.conf
exec "$powerdxx" --foreground > /dev/null
