# user/s6-rc/services/dwm/run
#!/bin/sh
cd "$HOME"
exec with-s6-env dwm
