# user/s6/.s6-svscan/finish
#!/usr/bin/env execlineb
importas scan S6_USER_SCAN
rm -rf ${scan}
